# crontab
*/1 * * * * /var/www/html php artisan schedule:run
